# deploy/podman/quadlet/booklib-db.container
[Unit]
Description=BookLib Database Container Using MariaDB

[Container]
ContainerName=booklib-db
Image=docker.io/library/mariadb:12.3.2
Pod=booklib.pod
Pull=always

# The official image's datadir. Deliberately a NEW volume: the old booklib-db
# volume holds a linuxserver /config tree, which this image would not recognise
# as a datadir. Mounting it here would initialise an empty database on top of it.
# See docs/UPGRADING-MARIADB-12.md before upgrading.
Volume=booklib-db-data:/var/lib/mysql

Environment=TZ=Etc/UTC
Environment=MYSQL_USER=booklib
Secret=booklib_db_pass,type=env,target=MYSQL_ROOT_PASSWORD
Secret=booklib_db_pass,type=env,target=MYSQL_PASSWORD
Environment=MYSQL_DATABASE=booklib

HealthCmd=CMD-SHELL healthcheck.sh --su-mysql --connect --innodb_initialized
HealthInterval=10s
HealthRetries=10
HealthStartPeriod=60s
HealthTimeout=5s

# Without this, systemd considers the unit started as soon as the container
# process exists, so booklib.container's After= would release the app while
# MariaDB is still initialising. Notify=healthy defers that until the
# healthcheck above passes.
Notify=healthy

DropCapability=all
# CHOWN/FOWNER: the entrypoint chowns the datadir to mysql before starting.
# DAC_OVERRIDE: root must traverse the mysql-owned 0700 datadir to do so.
# SETUID/SETGID: the entrypoint drops from root to mysql via gosu.
AddCapability=CAP_CHOWN CAP_DAC_OVERRIDE CAP_FOWNER CAP_SETGID CAP_SETUID
